# opensource/sysstat-6.0.1/sysstat.cron.hourly.in
#!/bin/sh
# Run system activity accounting tool every 10 minutes
SA_LIB_DIR/sa1 -d 600 6 &
